# files/supervisord-rabbitmq.conf
[program:rabbitmq]
command=/rabbitmq.sh
